# supervisorctl/site.conf
[program:ic3]
command=/home/ubuntu/ic3/venv/bin/uwsgi --ini uwsgi.ini 
directory=/home/ubuntu/ic3
;user = forge
stopsignal=QUIT
stdout_logfile=/tmp/ic3.stdout.log
stderr_logfile=/tmp/ic3.stderr.log